# src/recommendationservice/recomendations.service
[Unit]
Description=Payment Service

[Service]

WorkingDirectory=/home/botique/recommendation

Environment=PYTHONUNBUFFERED=1
Environment=ENABLE_PROFILER=1

ExecStart=/usr/bin/python3 recommendation_server.py
SyslogIdentifier=emailservice

[Install]
WantedBy=multi-user.target